# homeworks/homework1/systemd_files/check_status.timer
[Unit]
Description=Run Check Flask Web Application Status every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
